# pi/auto-update.service.example
# systemd USER service: one self-update check (run by the timer). The installer
# (install-kiosk.sh) writes this for you — this example is for reference / manual setup.
#   cp pi/auto-update.service.example ~/.config/systemd/user/dashboard-auto-update.service
#   cp pi/auto-update.timer.example   ~/.config/systemd/user/dashboard-auto-update.timer
#   systemctl --user daemon-reload && systemctl --user enable --now dashboard-auto-update.timer
# Watch it work:  journalctl --user -u dashboard-auto-update -f
#
# App-agnostic: auto-update.sh rebuilds whichever app changed and reloads the
# kiosk, which shows whatever kiosk.json selects. Adjust the path if you cloned
# somewhere other than ~/raspberry-playground.

[Unit]
Description=Dashboard self-update (git pull + build changed apps + reload kiosk)

[Service]
Type=oneshot
ExecStart=%h/raspberry-playground/pi/auto-update.sh
