# pi/dashboard.service.example
# systemd USER service for the dashboard server.
# Install: copy to ~/.config/systemd/user/worldcup-dashboard.service
#   systemctl --user daemon-reload
#   systemctl --user enable --now worldcup-dashboard.service
# Adjust WorkingDirectory if you cloned somewhere other than ~/world-cup-dashboard.
#
# Configuration (provider, port, API key) is read from the .env file in
# WorkingDirectory — see SETUP.md "Configure the data source". This unit
# deliberately sets no Environment= app vars, because those would override .env.

[Unit]
Description=World Cup Dashboard Server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/world-cup-dashboard
ExecStart=/usr/bin/npm run start
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
