# apps/family/pi/dashboard.service.example
# systemd USER service for the family dashboard server (port 3001).
# Install: copy to ~/.config/systemd/user/family-dashboard.service
#   systemctl --user daemon-reload
#   systemctl --user enable --now family-dashboard.service
# Adjust WorkingDirectory if you cloned somewhere other than ~/raspberry-playground.
#
# Configuration (provider, port, iCal URLs) is read from the .env file in
# WorkingDirectory. This unit deliberately sets no Environment= app vars,
# because those would override .env.

[Unit]
Description=Family Dashboard Server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/raspberry-playground/apps/family
ExecStart=/usr/bin/npm run start
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
